# sources/skel/.config/autostart/pulse-audio-server.desktop
[Desktop Entry]
Exec=/usr/bin/pulseaudio
Name=pulse-audio-server
Type=Application
Version=1.0
